# supervisor/viewer.conf
[program:viewer]
command=/usr/bin/pqiv --display=:0.0 media/images/
environment=PQIVRC_PATH="/home/pi/.config/pqivrc"
directory=/home/pi/shashinta/web
user=pi